# scripts/systemd/openclaw-n8n-stack.service
[Unit]
Description=OpenClaw n8n Stack (docker compose)
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/tommy/Dev/openclaw-n8n-blueprint
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
